# rootfs/etc/services.d/heimdall/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "Heimdall service finished"
